# docker/web/etc/cron.hourly/pay-subscription-fee
#!/bin/bash

cd /app/backend || exit
php yii payment/subscription-fee --interactive=0
